# deploy/funckcionario-pipeline.service
[Unit]
Description=funckcionario data pipeline
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=zafion
Group=zafion
WorkingDirectory=/home/zafion/Proyects/Radar_Docent_CV
EnvironmentFile=/etc/funckcionario/funckcionario.env
ExecStart=/home/zafion/Proyects/Radar_Docent_CV/scripts/run_pipeline.sh
