# deploy/funkcionario-centers.service
[Unit]
Description=Funkcionario sync catálogo de centros
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=oneshot
User=funkcionario
Group=funkcionario
WorkingDirectory=/srv/funkcionario/app
EnvironmentFile=/etc/funkcionario/funkcionario.env
ExecStart=/srv/funkcionario/app/.venv/bin/python run_sync_centers_catalog.py
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
