# deploy/funkcionario-web.service
[Unit]
Description=Funkcionario FastAPI web
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=funkcionario
Group=funkcionario
WorkingDirectory=/srv/funkcionario/app
EnvironmentFile=/etc/funkcionario/funkcionario.env
ExecStart=/srv/funkcionario/app/scripts/start_web.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
