# deploy/funkcionario-pipeline.service
[Unit]
Description=funkcionario data pipeline
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=zafion
Group=zafion
WorkingDirectory=/home/zafion/Proyects/Radar_Docent_CV
EnvironmentFile=/etc/funkcionario/funkcionario.env
ExecStart=/home/zafion/Proyects/Radar_Docent_CV/scripts/run_pipeline.sh
